# com.meetingtranscriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.meetingtranscriber</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/mantra/src/meeting-transcriber/.venv/bin/python</string>
        <string>/Users/mantra/src/meeting-transcriber/transcribe_meeting.py</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/mantra/MeetingInbox</string>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/meetingtranscriber.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/meetingtranscriber.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/mantra/src/meeting-transcriber</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
